# crawlers/crontab
0 15 * * * python3 /app/timepad.py
5 15 * * * python3 /app/meetup.py
