# backend/crawlers/crontab
00 15 * * * cd /app && python3 /app/timepad.py
05 15 * * * cd /app && python3 /app/meetup.py
10 15 * * * cd /app && python3 /app/devsday.py
